# gold-web.service
[Unit]
Description=Gold Price Analyzer Web Interface
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/gold-price-analyzer
Environment="PATH=/root/gold-price-analyzer/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/gold-price-analyzer/venv/bin/python /root/gold-price-analyzer/web_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target